# assets/systemd/rincron-mini.user.service
[Unit]
Description=Rincron-mini

[Service]
ExecStart=/usr/bin/rincron-mini
ExecReload=/bin/kill -USR1 $MAINPID
