# assets/deb/preinst
#! /bin/bash

systemctl stop rincron-mini
systemctl disable rincron-mini
